# live/engine/stock-signal-engine.service
[Unit]
Description=Pinch Stock Signal Engine
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/bob/AI_sandbox/pinch-stock-trading
ExecStart=/usr/bin/python3 live/engine/signal_engine.py
Environment=PINCH_DB=/mnt/media/market_data/pinch_market.db
Environment=FRED_API_KEY=3a0164644a3bbe2bb83d31c6b5209356
Environment=PATH=/home/bob/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
